# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:node-api]
command=node /app/src/index.js
directory=/app
environment=HTTP_PORT="80",NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:python-pay]
command=python3 -m uvicorn server:app --host 0.0.0.0 --port 8000
directory=/app
environment=PYTHONPATH="/usr/lib/python3.12/site-packages:/app"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
